# services/minq-ip-bans.service
[Unit]
Description=Minq ip bans
After=network.target

[Service]
Type=oneshot
ExecStart=minq-ip-bans

[Install]
WantedBy=multi-user.target
